# crontab
* * * * * php Laravel-Crypta/artisan schedule:run >> /dev/null 2>&1
